# packaging/vendi-chatd.service
[Unit]
Description=vendiMessage backend daemon (vendi-chatd)
Documentation=https://github.com/VendavalSC/vendiOS
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/vendi-chatd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
